# contrib/sprinkler-ai.service
[Unit]
Description=sprinkler-ai daily irrigation plan + execute
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/sprinkler-ai
ExecStart=%h/sprinkler-ai/.venv/bin/python -u -m sprinkler_ai.cli
StandardOutput=journal
StandardError=journal
